# addons/daemontools/service/a2b-callback-daemon-php/run
#!/bin/bash

sleep 5
exec 2>&1
date "+%F starting a2b-callback-daemon-php at %X"
exec /usr/local/src/a2billing/CallBack/callback-daemon-php/callback_daemon.php
